# pi/rc.local
#!/bin/sh -e

(sleep 4; python /opt/augment00/setup_wifi.py)&

exit 0